# packages/network/mirror/system.d/mirror.service
[Unit]
Description=Clue Mirror
After=network-online.target
Requires=mirror-defaults.service
ConditionPathExists=/home/.cache/services/mirror.conf

[Service]
Type=simple
ExecStart=/usr/bin/mirror --config-file=/home/.cache/services/mirror.conf
KillMode=process

[Install]
WantedBy=multi-user.target
