# packages/network/bluez/system.d/bluetooth-defaults.service
[Unit]
Description=Bluetooth defaults
After=local-fs.target
ConditionPathExists=!/home/.cache/services/bluez.conf
ConditionPathExists=!/home/.cache/services/bluez.disabled

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cp /usr/share/services/bluez.conf /home/.cache/services/'
RemainAfterExit=yes
