# packs/audio/alsa-utils/udev.d/90-alsa-restore.rules
# When a sound device is detected, restore the volume settings
  KERNEL=="controlC[0-9]*",  NAME="snd/%k", ACTION=="add", RUN+="soundconfig %k"
