# packages/system/busybox/system.d.opt/cron-defaults.service
[Unit]
Description=Cron defaults
After=local-fs.target

ConditionPathExists=!/home/.cache/services/crond.conf
ConditionPathExists=!/home/.cache/services/crond.disabled

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cp /usr/share/services/crond.conf /home/.cache/services/'
RemainAfterExit=yes
