# packs/addons/service/mariadb/source/profile.d/99-mysql-histfile.profile
export MYSQL_HISTFILE=/clue/.kodi/userdata/addon_data/service.mariadb/mysql_history
