# packages/system/busybox/system.d/rpi-flash-firmware.target
[Unit]
Description=RPi Flash Firmware target
Requires=boot-firmware.service
After=boot-firmware.service
AllowIsolate=yes
