# packages/system/busybox/system.d/storage-log.service
[Unit]
Description=Create Persistent Log Directory on /home
DefaultDependencies=no
RequiresMountsFor=/home

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/mkdir -p /home/.cache/log/journal
