# packages/system/busybox/system.d/var-log.mount
[Unit]
Description=Persistent Log Storage
RequiresMountsFor=/var /home
Requires=storage-log.service
After=storage-log.service
ConditionKernelCommandLine=!installer
ConditionKernelCommandLine=|debugging
ConditionPathExists=|/home/.cache/debug.clue

[Mount]
What=/home/.cache/log
Where=/var/log
Options=bind
LazyUnmount=yes
